# fangofmongo.service.conf
[program:mongo]
directory=/Fang-of-Mongo/fangofmongo
command=/usr/bin/python manage.py runserver 0.0.0.0:8000
