# mongoexpress.service.conf
[program:mongo]
directory=/andzdroid-mongo-express-*
command=/usr/bin/npm start
