# mviewer.service.conf
[program:mongo]
directory=/mViewer
command=/bin/bash /mViewer/start_mViewer.sh
